# webserver/supervisor.sample.conf
[program:djangopro]   ##修改成你的项目名称
command=/path/to/start.sh
environment=DJANGO_ENV="PRO"
user=root  ##修改成你的用户
autostart=true
autorestart=true
redirect_stderr=True