# webserver/supervisor.sample.conf
[program:PKU_PHY_SU]   ##修改成你的项目名称
command=/path/to/PKU_PHY_SU/webserver/start.sh
environment=DJANGO_ENV="PRO"
user=root  ##修改成你的用户
autostart=true
autorestart=true
redirect_stderr=True